# supervisord.conf
[program:wechat]
directory=/root/wechat/
command=uwsgi --ini uwsgi_weixin.ini

stdout_logfile=/data/logs/uwsgi_wechat.log

autostart=true
autorestart=true
redirect_stderr=true

stopsignal=QUIT
